# sources/meta-syncroness/recipes-core/init-lowlevel/files/init-lowlevel.service
[Unit]
Description=Low Level Hardware Init
After=basic.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/etc/init-lowlevel.sh