# sources/meta-cytovale/recipes-core/cym-app-service/cym-app-service/cym-app.service
[Unit]
Description=Cytovale Application Service
After=nvm-partition.service

[Service]
Type=notify
TimeoutStartSec=300
ExecStart=/bin/bash -c "/home/root/cytovale-sw/scripts/start.sh"
ExecStop=/bin/bash -c "/home/root/cytovale-sw/scripts/stop.sh"
NotifyAccess=all
WatchdogSec=60
Restart=always
StartLimitInterval=180min
StartLimitBurst=1
StartLimitAction=reboot-immediate

[Install]
WantedBy=multi-user.target