# sources/meta-cytovale/recipes-core/cym-webserver/cym-webserver/cym-webserver.service
[Unit]
Description=Cytovale Django Webserver
After=multi-user.target

[Service]
Type=notify
ExecStart=/bin/bash -c "/home/root/cym-webserver.sh"

[Install]
After=multi-user.target
WantedBy=multi-user.target
